# etc/update-motd.d/01-uptime
#! /bin/bash
echo "### System uptime ###"
uptime=$(uptime --pretty)
echo -e System.uptime=\"${uptime#up }\"
echo
