# etc/update-motd.d/02-date
#! /bin/bash
echo "### UTC date/time ###"
currentime=$(date --utc +%D\ %T)
echo -e System.Time=\"${currentime}\"
echo
