# etc/update-motd.d/06-logs
#! /bin/bash
echo "### Latest entries in auth log ###"
tail -n5 /var/log/auth.log |sed s/"]:"/"\n"/ | grep -v $(date +%b)
echo
